# gce/systemd/billing-discord-notifier.service
[Unit]
Description=Billing Budget Alert → Discord Notifier
Documentation=https://github.com/Tagomori0211/Minecraft_java_k3s
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
# --network=host: GCE インスタンスメタデータ (169.254.169.254) に到達するために必要
ExecStart=/usr/bin/docker run --rm \
  --network=host \
  -v /opt/mc-proxy/gce/scripts:/scripts:ro \
  google/cloud-sdk:slim \
  python3 /scripts/billing-discord-notifier.py
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal
